# system_files/nix/usr/lib/systemd/system/nix-relabel.timer
[Unit]
Description=Run Nix SELinux relabel weekly

[Timer]
OnCalendar=weekly
RandomizedDelaySec=30m
Persistent=true
Unit=nix-relabel.service

[Install]
WantedBy=timers.target
